# bots/outreach/systemd/kitaiservice-outreach-daemon@.service
[Unit]
Description=KitaiService Outreach Daemon (group %i)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kitaiservice_outreach
EnvironmentFile=/opt/kitaiservice_outreach/.env
ExecStart=/opt/kitaiservice_outreach/venv/bin/python3 -u manager_daemon.py --group %i
Restart=on-failure
RestartSec=10
TimeoutStopSec=30
KillMode=mixed
StandardOutput=append:/opt/kitaiservice_outreach/daemon.log
StandardError=append:/opt/kitaiservice_outreach/daemon.log

[Install]
WantedBy=multi-user.target
